# systemd/sourcecatcher-update.timer
[Unit]
Description=Update Sourcecatcher database

[Timer]
OnCalendar=*-*-* 12:00:00
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
